# misc/systemd/user/org.desktopspec.ApplicationManager1.service.in
# SPDX-FileCopyrightText: 2023 UnionTech Software Technology Co., Ltd.
#
# SPDX-License-Identifier: LGPL-3.0-or-later

[Unit]
Description=Deepin Application Manager

[Service]
Type=dbus
BusName=org.desktopspec.ApplicationManager1
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/dde-application-manager-reborn
Restart=always
Environment=QT_LOGGING_RULES="*.debug=false"
PrivateUsers=false
# turn off PrivateUser to prevent AM can't access some directory. eg. "/persistent/linglong"
