# misc/systemd/user/org.desktopspec.ApplicationManager1.service.in
# SPDX-FileCopyrightText: 2023 - 2026 UnionTech Software Technology Co., Ltd.
#
# SPDX-License-Identifier: LGPL-3.0-or-later

[Unit]
Description=Deepin Application Manager
StartLimitBurst=3
CollectMode=inactive-or-failed

PartOf=dde-session-pre.target
Before=dde-session-pre.target

[Service]
Type=dbus
BusName=org.desktopspec.ApplicationManager1
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/dde-application-manager
Environment=QT_LOGGING_RULES="*.debug=false"
Slice=session.slice
Restart=always
RestartSec=500ms
